# root/etc/services.d/rclone_mounter/run
#!/usr/bin/with-contenv bash
set -x


_term() {
  echo "Caught SIGTERM signal!"
  wg-quick down wg0
}

trap _term SIGTERM

trap _term SIGTERM

echo "running"
env


sleep infinity &

wait